# php8.0/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/usr/bin/with-contenv sh
/usr/sbin/php-fpm8.0 --nodaemonize
