# etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
su appuser
exec /usr/sbin/php-fpm8.0 --nodaemonize
