# etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
/usr/sbin/php-fpm8.0 --nodaemonize
